# episode-3/ConnectionInfoApi/systemd/connection-info-api.service
[Unit]
Description=Connection Info API
After=network-online.target
Wants=network-online.target

[Service]
Type=notify

# Least-privilege user and group
User=connectioninfo
Group=connectioninfo

# App directory
WorkingDirectory=/opt/connection-info-api

# Start the ASP.NET Core self-contained single-file app
ExecStart=/opt/connection-info-api/ConnectionInfoApi --urls=http://0.0.0.0:57171

# Restart on failure
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

# Optional: post-start health check for observability
ExecStartPost=/usr/bin/curl --fail --silent --max-time 5 http://localhost:57171/healthz

# Runtime environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT=true

# Logging (systemd journal)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=connection-info-api

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
PrivateDevices=true
RestrictNamespaces=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=~

# Network family restrictions (IPv4, IPv6, UNIX sockets)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Runtime directory if used
RuntimeDirectory=connectioninfo
RuntimeDirectoryMode=0750

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
